# etc/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:coolsite]
command=/home/sites/coolsite/env/bin/uwsgi /home/sites/coolsite/wsgi.ini
stdout_logfile=/home/sites/coolsite/logs/wsgi.log
stderr_logfile=/home/sites/coolsite/logs/wsgi_err.log
autostart=true
autorestart=true
redirect_stderr=true

[program:coolsite_celery]
command=/home/sites/coolsite/project/manage.py celeryd --loglevel=INFO -B
stdout_logfile=/home/sites/coolsite/logs/celery.log
stderr_logfile=/home/sites/coolsite/logs/celery_err.log
redirect_stderr=true
autostart=true
autorestart=true
startsecs=5
; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 60
